# prometheus/service/redis_exporter.service
[Unit]
Description=Redis Exporter
Wants=network-online.target
After=network-online.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/redis_exporter -web.listen-address ":9121" -redis.addr "redis://0.0.0.0:6379" -redis.password "YourPassword"

[Install]
WantedBy=multi-user.target